# backup/systemd/user/tg-autorestart.service
[Unit]
Description=restart telegram-desktop after crash

[Service]
ExecStart=/usr/bin/telegram-desktop
Restart=always