# backup/systemd/user/v2raya.service
[Unit]
Description=v2raya

[Service]
ExecStart=/usr/bin/v2raya --lite
Restart=always
