# backup/systemd/user/v2raya.service
[Unit]
Description=v2raya: the v2ray client

[Service]
ExecStart=/usr/bin/v2raya --lite
Restart=always
